# deploy/systemd/agentic-historian.service.d/restart.conf
# P3-2 (#247): keep the bot alive across crashes AND across the /update self-exit.
# Install on tei:  sudo cp this file to
#   /etc/systemd/system/agentic-historian.service.d/restart.conf
# then:  sudo systemctl daemon-reload && sudo systemctl restart agentic-historian
# Restart=always restarts the unit on any exit (incl. exit 0 from /update), so
# systemd brings the bot back on the new code after it exits.
[Service]
Restart=always
RestartSec=2
